# launchpad-init.service
[Unit]
Description=Run launchpad-init

[Service]
ExecStart=%h/.local/bin/launchpad-init
Type=simple

[Install]
WantedBy=default.target
